# etc/init.d/sudo-openrc
#!/sbin/openrc-run

description="Initializes the sudo timestamp directory"

depend() {
	need localmount
	provide sudo
}

start() {
	ebegin "Initializing sudo timestamp directory"

	checkpath -d -m 0711 -o root:root /run/sudo
	checkpath -d -m 0700 -o root:root /run/sudo/ts

	find /run/sudo/ts -type f -exec touch -d @0 '{}' +

	[ -x /usr/sbin/restorecon ] && /usr/sbin/restorecon /run/sudo /run/sudo/ts

	eend 0
}
